# check_server.service
[Unit]
Description=Check Server Health

[Service]
ExecStart=/home/arboratorgrew/arborator-parser/check_if_server_is_functionning.sh
